# systemd-service/httpshare.service
[Unit]
Description=HTTPShare Image Server
After=network.target
Wants=network.target

[Service]
Type=simple
User=mon
Group=mon
WorkingDirectory=/home/mon/httpshare
Environment=DISPLAY=:0
Environment=XAUTHORITY=/home/mon/.Xauthority
ExecStart=/usr/bin/node /home/mon/httpshare/server.js
ExecStop=/home/mon/httpshare/stop-server.sh
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target 